# deploy/deploy-package/draw-and-guess-python.service
[Unit]
Description=Draw and Guess Game Server (Python)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/draw-and-guess
ExecStart=/opt/draw-and-guess/venv/bin/uvicorn main:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5
Environment=STATIC_DIR=client/dist
Environment=PORT=3000
Environment=DB_PATH=/opt/draw-and-guess/game.db
Environment=LOG_LEVEL=info
Environment=PUBLIC_URL=
Environment=ADMIN_PASSWORD=

[Install]
WantedBy=multi-user.target
